# systemd/disk-maintenance-report.service
[Unit]
Description=disk-maintenance (learning/report)

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/disk-maintenance report
