# systemd/disk-maintenance-emergency.service
[Unit]
Description=disk-maintenance (emergency)

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/disk-maintenance emergency
