# systemd/disk-maintenance-cleanup.service
[Unit]
Description=disk-maintenance (guarded cleanup)

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/disk-maintenance cleanup
