# systemd/disk-maintenance-report.timer
[Unit]
Description=disk-maintenance report timer

[Timer]
OnCalendar=hourly
RandomizedDelaySec=10m
Persistent=true

[Install]
WantedBy=timers.target
